# basic-startup-conf/data/etc/systemd/system/pretty-early.target
[Unit]
Description=Pretty Early Application
Before=basic.target
